# infra/deploy/warehouse-platform-backup.service
[Unit]
Description=warehouse-platform PostgreSQL backup

[Service]
Type=oneshot
WorkingDirectory=/var/www/warehouse-platform
EnvironmentFile=/var/www/warehouse-platform/.env
ExecStart=/usr/bin/npm run backup:postgres
