# infra/deploy/warehouse-platform-backup.timer
[Unit]
Description=Run warehouse-platform PostgreSQL backup nightly

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true

[Install]
WantedBy=timers.target
